# Bash_scripts/HW21/in_hw20/roles/hw202/templates/logsysd.timer.j2
[Unit]
Description=timer for system_monitor

[Timer]
OnUnitActiveSec=30s
OnBootSec=10s

[Install]
WantedBy=timers.target
